# provision/files/etc/xdg/autostart/kiwix.desktop
[Desktop Entry]
Name=kiwix
Type=Application
Encoding=UTF-8
Version=1.0
Exec=/home/vagrant/kiwix/kiwix
X-GNOME-Autostart-enabled=true